# supervisord.conf
[program:my_server]
command=python3 index.py
autostart=true
autorestart=true
stderr_logfile=/path/to/logs/stderr.log
stdout_logfile=/path/to/logs/stdout.log
